# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g 'daemon off;'
priority=5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ranker]
command=bash -c "set -a; . /opt/mastodon/.env.runtime; set +a; exec /opt/rankers/.venv/bin/python -m uvicorn corda_rankers.app:app --host 127.0.0.1 --port ${CORDA_RANKER_PORT}"
directory=/opt/rankers
user=mastodon
environment=HOME="/opt/mastodon",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
priority=15
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=bash -c "set -a; . /opt/mastodon/.env.runtime; set +a; exec bundle exec puma -C config/puma.rb"
directory=/opt/mastodon
user=mastodon
environment=HOME="/opt/mastodon",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",PORT="3000"
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
